# packaging/postinst
#!/bin/sh
launchctl load /Library/LaunchDaemons/jetslammed.plist 2>/dev/null

